# systemd-user/walker.service
[Unit]
Description=Walker application launcher (service mode)
After=graphical-session.target elephant.service
Requires=elephant.service

[Service]
ExecStart=/usr/bin/walker --gapplication-service
Restart=always
RestartSec=2
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
